# deploy/news-troubleshoot.timer
# News App - Troubleshoot Timer
#
# Runs daily at 07:00 to check for failed jobs.
#
[Unit]
Description=Run news-app troubleshoot check daily

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
