# deploy/news-troubleshoot.service
# News App - Troubleshoot Service
#
# Automatically diagnoses failed job runs by querying the database
# and creating a Shelley conversation with diagnostic information.
#
# See docs/DEPLOYMENT.md for details.
#
[Unit]
Description=Troubleshoot news-app job runs
After=network.target

[Service]
Type=oneshot
ExecStart=/home/exedev/news-app/news-app troubleshoot
User=exedev
WorkingDirectory=/home/exedev/news-app

[Install]
WantedBy=multi-user.target
