# deploy/news-cleanup.service
# News App - Conversation Cleanup Service
#
# Cleans up old Shelley conversations (older than 48 hours)
# to prevent database bloat.
#
# See docs/DEPLOYMENT.md for details.
#
[Unit]
Description=Cleanup old news-app conversations
After=network.target

[Service]
Type=oneshot
ExecStart=/home/exedev/news-app/news-app cleanup
User=exedev
WorkingDirectory=/home/exedev/news-app

[Install]
WantedBy=multi-user.target
